# ops/systemd/auto-trader-notify.service.example
[Unit]
Description=Auto Trader Notification Watcher
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/komug/projects/auto_trader
Environment=AUTO_TRADER_NOTIFY_WEBHOOK_URL=https://example.com/webhook
ExecStart=/home/komug/projects/auto_trader/.venv/bin/python -m auto_trader.notify --from-env --watch --interval-sec 5 --output-dir data/ops
Restart=always
RestartSec=3
User=komug

[Install]
WantedBy=multi-user.target
