# ops/systemd/auto-trader-weekly-revalidation.service.example
[Unit]
Description=Auto Trader Weekly Revalidation
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/komug/projects/auto_trader
Environment=PYTHONUNBUFFERED=1
Environment=RUN_ROOT=data/validation/weekly_autotune
ExecStart=/home/komug/projects/auto_trader/scripts/weekly_autotune_pipeline.sh
StandardOutput=journal
StandardError=journal
