# ops/systemd/auto-trader-ohlcv-refresh.user.service.example
[Unit]
Description=Auto Trader OHLCV Refresh for Positions (User)
After=auto-trader-worker.service

[Service]
Type=oneshot
WorkingDirectory=/home/komug/projects/auto_trader
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/komug/projects/auto_trader/scripts/refresh_ohlcv_for_positions.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
